# addons/claude-backend/rootfs/etc/s6-overlay/s6-rc.d/claude-backend/finish
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Claude Backend finish script
# ==============================================================================
if test "$1" -eq 256 ; then
  e=$((128 + $2))
else
  e="$1"
fi

if test "$e" -ne 0; then
    bashio::log.warning "claude-backend exited with code $e (by signal $2)."
else
    bashio::log.info "claude-backend exited cleanly (code $e)."
fi

# Always halt the container — prevents s6 from restarting the service
# and ensures the Supervisor sees the addon as properly stopped.
echo "$e" > /run/s6-linux-init-container-results/exitcode
exec /run/s6/basedir/bin/halt
