# addons/claude-backend/rootfs/etc/s6-overlay/s6-rc.d/claude-backend/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start AI Assistant Backend API
# ==============================================================================
bashio::log.info "Starting AI Assistant Backend API..."

# AI Provider configuration
AI_PROVIDER=$(bashio::config 'ai_provider')
ANTHROPIC_API_KEY=$(bashio::config 'anthropic_api_key')
OPENAI_API_KEY=$(bashio::config 'openai_api_key')
GOOGLE_API_KEY=$(bashio::config 'google_api_key')
GITHUB_TOKEN=$(bashio::config 'github_token')
GITHUB_MODEL=$(bashio::config "github_model" 2>/dev/null || echo "")
if [ "${GITHUB_MODEL}" = "null" ]; then GITHUB_MODEL=""; fi
API_PORT=$(bashio::config 'api_port')
DEBUG_MODE=$(bashio::config 'debug_mode')
ENABLE_FILE_ACCESS=$(bashio::config 'enable_file_access')

# SUPERVISOR_TOKEN: injected by HA Supervisor, fallback via bashio
if [[ -z "${SUPERVISOR_TOKEN:-}" ]]; then
    SUPERVISOR_TOKEN="${HASSIO_TOKEN:-}"
fi
if [[ -z "${SUPERVISOR_TOKEN:-}" ]]; then
    SUPERVISOR_TOKEN=$(bashio::api.supervisor_token 2>/dev/null || true)
fi

export AI_PROVIDER
export ANTHROPIC_API_KEY
export OPENAI_API_KEY
export GOOGLE_API_KEY
export GITHUB_TOKEN
export GITHUB_MODEL
export API_PORT
export DEBUG_MODE
export ENABLE_FILE_ACCESS
export SUPERVISOR_TOKEN

bashio::log.info "AI Provider: ${AI_PROVIDER}"
bashio::log.info "SUPERVISOR_TOKEN present: $([ -n \"${SUPERVISOR_TOKEN}\" ] && echo 'yes' || echo 'NO - HA API will fail!')"
bashio::log.info "SUPERVISOR_TOKEN length: ${#SUPERVISOR_TOKEN}"

cd /app
exec python3 api.py
