# addons/claude-backend/rootfs/etc/s6-overlay/s6-rc.d/claude-backend/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start AI Assistant Backend API
# ==============================================================================
bashio::log.info "Starting AI Assistant Backend API..."

# Startup configuration (from config.yaml)
export DEBUG_MODE=$(bashio::config 'debug_mode')
export COLORED_LOGS=$(bashio::config 'colored_logs')
export LOG_LEVEL=$(bashio::config 'log_level')

# API keys (from config.yaml)
export ANTHROPIC_API_KEY=$(bashio::config 'anthropic_api_key')
export OPENAI_API_KEY=$(bashio::config 'openai_api_key')
export GOOGLE_API_KEY=$(bashio::config 'google_api_key')
export NVIDIA_API_KEY=$(bashio::config 'nvidia_api_key')
export GITHUB_TOKEN=$(bashio::config 'github_token')
export GROQ_API_KEY=$(bashio::config 'groq_api_key')
export MISTRAL_API_KEY=$(bashio::config 'mistral_api_key')
export OPENROUTER_API_KEY=$(bashio::config 'openrouter_api_key')
export DEEPSEEK_API_KEY=$(bashio::config 'deepseek_api_key')
export XAI_API_KEY=$(bashio::config 'xai_api_key')
export MINIMAX_API_KEY=$(bashio::config 'minimax_api_key')
export AIHUBMIX_API_KEY=$(bashio::config 'aihubmix_api_key')
export SILICONFLOW_API_KEY=$(bashio::config 'siliconflow_api_key')
export VOLCENGINE_API_KEY=$(bashio::config 'volcengine_api_key')
export DASHSCOPE_API_KEY=$(bashio::config 'dashscope_api_key')
export MOONSHOT_API_KEY=$(bashio::config 'moonshot_api_key')
export ZHIPU_API_KEY=$(bashio::config 'zhipu_api_key')
export PERPLEXITY_API_KEY=$(bashio::config 'perplexity_api_key')
export CUSTOM_API_KEY=$(bashio::config 'custom_api_key')
export CUSTOM_API_BASE=$(bashio::config 'custom_api_base')
export CUSTOM_MODEL_NAME=$(bashio::config 'custom_model_name')
export OLLAMA_BASE_URL=$(bashio::config 'ollama_base_url')
export OLLAMA_API_KEY=$(bashio::config 'ollama_api_key')
# OAuth-only providers
export GITHUB_COPILOT_TOKEN=""
export OPENAI_CODEX_TOKEN=""

# SUPERVISOR_TOKEN: injected by HA Supervisor, fallback via bashio
if [[ -z "${SUPERVISOR_TOKEN:-}" ]]; then
    SUPERVISOR_TOKEN="${HASSIO_TOKEN:-}"
fi
if [[ -z "${SUPERVISOR_TOKEN:-}" ]]; then
    SUPERVISOR_TOKEN=$(bashio::api.supervisor_token 2>/dev/null || true)
fi
export SUPERVISOR_TOKEN

bashio::log.info "DEBUG_MODE: ${DEBUG_MODE}"
bashio::log.info "LOG_LEVEL: ${LOG_LEVEL}"
bashio::log.info "SUPERVISOR_TOKEN: $([ -n "${SUPERVISOR_TOKEN}" ] && echo 'present' || echo 'MISSING')"
bashio::log.info "Runtime settings managed via /config/amira/settings.json"

# Remove macOS AppleDouble metadata files (._*) that Finder creates on SMB/AFP shares
find /app/translations -name '._*' -delete 2>/dev/null || true

cd /app

exec python3 server.py
